# recipes-qt5/packages/qt5-demo/qt5-demo.service
# See http://www.dsm.fordham.edu/cgi-bin/man-cgi.pl?topic=systemd.service

[Unit]
Description=Qt5 Demo application
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/qt5-demo
RemainAfterExit=no

[Install]
WantedBy=multi-user.target